# utils/systemd/system/batch-batch.slice
## live-net-initramfs - System Configuration Components
##
## live-net-initramfs contains the components that configure a live system during
## the boot process (late userspace).
##

[Unit]
Description=Batch jobs slice.
Documentation=man:live-net-initramfs
DefaultDependencies=no
Before=slices.target
Wants=-.slice
After=-.slice
Wants=batch.slice
After=batch.slice

[Slice]
CPUWeight=400
X-ControlGroupAttribute=cpu.policy:batch
