# utils/systemd/system/interactive-intlp.slice
## live-net-initramfs - System Configuration Components
##
## live-net-initramfs contains the components that configure a live system during
## the boot process (late userspace).
##

[Unit]
Description=Interactive low priority tasks slice.
Documentation=man:live-net-initramfs
DefaultDependencies=no
Before=slices.target
Wants=-.slice
After=-.slice
After=interactive.slice
Wants=interactive.slice

[Slice]
CPUWeight=400
X-ControlGroupAttribute=cpu.policy:batch
