# utils/systemd/system/system-ctl.slice
## live-net-initramfs - System Configuration Components
##
## live-net-initramfs contains the components that configure a live system during
## the boot process (late userspace).
##

[Unit]
Description=system control tasks slice.
Documentation=man:live-net-initramfs
DefaultDependencies=no
Before=slices.target
Wants=-.slice
After=-.slice
After=system.slice
Wants=system.slice

[Slice]
CPUWeight=900
X-ControlGroupAttribute=cpu.nice:-10
